# packages/testing/kubernetes/kube-proxy/kube-proxy.initd
#!/sbin/openrc-run

NAME=kube-proxy
DAEMON=/usr/bin/$NAME
OPTS="
        --cluster-cidr=${CLUSTER_CIDR} \
        --kubeconfig=${KUBECONFIG} \
        --proxy-mode=${PROXY_MODE} \
        --v=${LOG_LEVEL}"

depend() {
        need net
        after firewall
}

start() {
        ebegin "Starting ${NAME}"
                start-stop-daemon --start \
                        --background \
                        --make-pidfile --pidfile /var/run/$NAME.pid \
                        --stderr "/var/log/$NAME.log" \
                        --stdout "/var/log/$NAME.log" \
                        --exec $DAEMON \
                        -- $OPTS
        eend $?
}

stop () {
        ebegin "Stopping ${NAME}"
                start-stop-daemon --stop \
                        --pidfile /var/run/$NAME.pid \
                        --exec $DAEMON
        eend $?
}