# packages/testing/kubernetes/kubelet/kubelet.initd
#!/sbin/openrc-run

NAME=kubelet
DAEMON=/usr/bin/$NAME

OPTS="
  --cloud-provider= \
  --container-runtime=${CONTAINER_RUNTIME} \
  --container-runtime-endpoint=${CONTAINER_RUNTIME_ENDPOINT} \
  --image-pull-progress-deadline=${IMAGE_PULL_PROGRESS_DEADLINE} \
  --kubeconfig=${KUBECONFIG} \
  --network-plugin=${NETWORK_PLUGIN} \
  --register-node=${REGISTER_NODE} \
  --v=${LOG_LEVEL}"

depend() {
        need net
        after firewall
}

start() {
        ebegin "Starting ${NAME}"
                start-stop-daemon --start \
                        --background \
                        --make-pidfile --pidfile /var/run/$NAME.pid \
                        --stderr "/var/log/$NAME.log" \
                        --stdout "/var/log/$NAME.log" \
                        --exec $DAEMON \
                        -- $OPTS
        eend $?
}

stop () {
        ebegin "Stopping ${NAME}"
                start-stop-daemon --stop \
                        --pidfile /var/run/$NAME.pid \
                        --exec $DAEMON
        eend $?
}