# packages/node_exporter/node_exporter.post-install
#!/bin/sh

rc-update add node_exporter

exit 0